# packages/prometheus-disk-exporter/prometheus-disk-exporter-0.0.1.tar.gz/prometheus-disk-exporter-0.0.1/systemd/prometheus-disk-exporter.service
[Unit]
Description=Service to export S.M.A.R.T. and partition data as Prometheus metrics

[Service]
ExecStart=/usr/local/bin/prometheus-disk-exporter
Type=simple
User=prometheus
Group=prometheus

[Install]
WantedBy=multi-user.target
